# mysql-backup.timer
[Unit]
Description=Run MySQL backup daily
Requires=mysql-backup.service

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target